# deploy/systemd/trener-rsvp-reminders.timer
[Unit]
Description=trener — ежедневный запуск RSVP-напоминаний (09:00 MSK)

[Timer]
# Бокс — Europe/Moscow; суффикс зоны пинит время даже при смене TZ бокса (systemd>=252; на боксе 255).
# Окно эндпоинта — ближайшие 48ч, так что нереспондент получит ~2 напоминания до тренировки.
OnCalendar=*-*-* 09:00:00 Europe/Moscow
Persistent=true        # пропущенный запуск (бокс был выключен) — выполнить при старте
RandomizedDelaySec=3m  # не бить ровно в 09:00 одновременно с соседними проектными таймерами
AccuracySec=1m

[Install]
WantedBy=timers.target
