# deploy/systemd/trener-backup.timer
[Unit]
Description=trener — ежедневный offsite-бэкап БД (03:30 MSK)

[Timer]
# Off-peak; бокс в Europe/Moscow (зону пиним явно). Бэкап не time-critical → широкий джиттер.
OnCalendar=*-*-* 03:30:00 Europe/Moscow
Persistent=true         # пропущенный запуск (бокс был выключен) — выполнить при старте
RandomizedDelaySec=10m
AccuracySec=1m

[Install]
WantedBy=timers.target
